# files/etc/crontabs/root
#publish map data
*/5 * * * * /usr/sbin/print_map.sh -p
